# etc/systemd/system/at-archiver.service
[Unit]
Description=ATArchiver Services

[Service]
Type=oneshot
ExecStart=/bin/true
RemainAfterExit=yes

[Install]
WantedBy=at-archiver.target
